# py-xdesire/supervisord.conf
[program:proxy]
command = python /home/dotcloud/current/proxy.py
